# private_dot_config/systemd/user/vdirsyncer.service
[Unit]
Description=vdirsyncer synchronization

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
